# release/supervisor.example.ini
[program:follow-monitor-api]
command=/www/wwwroot/follow-monitor/follow-monitor/start.sh
directory=/www/wwwroot/follow-monitor/follow-monitor/
autorestart=true
startsecs=3
startretries=3
stdout_logfile=/www/server/panel/plugin/supervisor/log/follow-monitor-api.out.log
stderr_logfile=/www/server/panel/plugin/supervisor/log/follow-monitor-api.err.log
stdout_logfile_maxbytes=2MB
stderr_logfile_maxbytes=2MB
user=root
priority=999
numprocs=1
minfds=65535

[program:follow-monitor-worker]
command=/www/wwwroot/follow-monitor/follow-monitor/start-worker.sh
directory=/www/wwwroot/follow-monitor/follow-monitor/
autorestart=true
startsecs=3
startretries=3
stdout_logfile=/www/server/panel/plugin/supervisor/log/follow-monitor-worker.out.log
stderr_logfile=/www/server/panel/plugin/supervisor/log/follow-monitor-worker.err.log
stdout_logfile_maxbytes=2MB
stderr_logfile_maxbytes=2MB
user=root
priority=998
numprocs=1
minfds=65535
